# systemd/minecraft-netbird-forward.service
[Unit]
Description=NetBird-private Minecraft and voice-chat forwards to VM 110
After=network-online.target netbird.service pve-firewall.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/local/sbin/minecraft-netbird-forward start
ExecStop=/bin/bash /usr/local/sbin/minecraft-netbird-forward stop

[Install]
WantedBy=multi-user.target
