# systemd/nazar-disk-usage-check.service
[Unit]
Description=Nazar disk usage alert check
After=pve-cluster.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nazar-disk-usage-check
