# systemd/nazar-disk-usage-check.timer
[Unit]
Description=Run Nazar disk usage alert check hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
