# systemd/nazar-zellij-web.service
[Unit]
Description=Nazar private Zellij web terminal
Documentation=https://zellij.dev/documentation/web-client.html
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alex
Group=alex
WorkingDirectory=/home/alex
Environment=HOME=/home/alex
Environment=USER=alex
Environment=LOGNAME=alex
ExecStart=/nix/var/nix/profiles/default/bin/zellij web --ip 127.0.0.1 --port 8082
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/alex

[Install]
WantedBy=multi-user.target
