# systemd/nazar-boot-alert.service
[Unit]
Description=Nazar boot/reboot alert
After=pve-cluster.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nazar-boot-alert

[Install]
WantedBy=multi-user.target
